# epm.list
%product Kung Foo Firewall
%copyright 1999-2005 by Foo Industries, All Rights Reserved.
%vendor Foo Industries
%license COPYING
%readme README
%description Kung Foo firewall software for your firewall.
%version 1.2.3p4 1020304

type mode owner group destination source options

d 755 root sys /var/log/npa/ -
d 755 root sys /usr/lib -
f 755 root sys /usr/bin/ bin/*
c 644 root sys /etc/npa/ etc/*
l 000 root sys /etc/init.d/npa /usr/bin/npa
f 0444 root sys /usr/share/doc/foo *.html
